# automation/launchd/com.personal-tracker.telegram-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Telegram Polling LaunchAgent for Personal Tracker

  Runs the Telegram poll script every 15 minutes to check for new
  Apple Health JSON exports sent via Telegram bot.

  Install:
    cp automation/launchd/com.personal-tracker.telegram-poll.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.personal-tracker.telegram-poll.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.personal-tracker.telegram-poll.plist
    rm ~/Library/LaunchAgents/com.personal-tracker.telegram-poll.plist

  Check status:
    launchctl list | grep personal-tracker

  View logs:
    tail -f ~/Documents/personal-tracker/data/logs/telegram_poll.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.personal-tracker.telegram-poll</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ethan.eismann/Documents/personal-tracker/.venv/bin/python</string>
        <string>/Users/ethan.eismann/Documents/personal-tracker/automation/telegram_poll.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ethan.eismann/Documents/personal-tracker</string>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/ethan.eismann/Documents/personal-tracker/data/logs/telegram_poll.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ethan.eismann/Documents/personal-tracker/data/logs/telegram_poll.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
