# automation/launchd/com.personal-tracker.telegram-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template only.

  Use this command on each machine to generate/install a machine-specific plist:
    automation/launchd/install_telegram_poll_launchagent.sh
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.personal-tracker.telegram-poll</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_ROOT__/.venv/bin/python</string>
    <string>__PROJECT_ROOT__/automation/telegram_poll.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_ROOT__</string>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>StandardOutPath</key>
  <string>__PROJECT_ROOT__/data/logs/telegram_poll.log</string>

  <key>StandardErrorPath</key>
  <string>__PROJECT_ROOT__/data/logs/telegram_poll.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
